# root/etc/s6-overlay/s6-overlay/s6-rc.d/ffdl/run
#!/command/with-contenv bash

exec s6-setuidgid abc /app/run.sh;