# root/etc/services.d/ffdl/run
#!/usr/bin/with-contenv bash

sleep 30s
exec s6-setuidgid abc /app/run.sh